# infra/systemd/castle-watchdog.timer
[Unit]
Description=Run the castle watchdog hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
