# infra/systemd/purge-auth-logs.service
# Daily retention purge: deletes login/audit log rows older than 90 days
# from the notulen database (see infra/nextcloud/purge_old_auth_logs.py).
# Data minimisation: logs you no longer need are logs that cannot leak.
[Unit]
Description=Purge auth/audit log rows older than 90 days
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
User=castle
Group=castle
WorkingDirectory=/opt/castle/repo
ExecStart=/usr/bin/python3 /opt/castle/repo/infra/nextcloud/purge_old_auth_logs.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
